# reboot_server1.service
[Unit]
Description=Server 1 Auto Reboot Service
After=network.target docker.target
Wants=docker.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/server-autoreboot/reboot_service.py /opt/server-autoreboot/configs/server1.env
User=root
Group=root
StandardOutput=journal
StandardError=journal
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target